# dist/packages/ovs-exporter/config/ovs-exporter.service
[Unit]
Description=OVS metrics exporter service

[Service]
ExecStart=/opt/asgard/ovs-exporter/ovn-kube-util --loglevel=4 ovs-exporter --metrics-bind-address 0.0.0.0:9310

[Install]
WantedBy=multi-user.target